# bridge/camera-bridge.service
[Unit]
Description=Surface Laptop 2 OV9734 webcam bridge
After=wireplumber.service
Wants=wireplumber.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/share/surface-laptop-2-camera/camera-bridge.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
