# bridge/camera-bridge.service
[Unit]
Description=Surface Laptop 2 Camera Bridge (ov9734 → PipeWire + V4L2)
After=pipewire.service
Requires=pipewire.service

[Service]
Type=simple
ExecStart=%h/.local/bin/camera-bridge
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
